# deploy/shred.service
[Unit]
Description=shred — zero-knowledge ephemeral file sharing
After=network.target

[Service]
Type=simple
User=shred
Group=shred
WorkingDirectory=/opt/shred
EnvironmentFile=/opt/shred/.env

# PrivateTmp's default /tmp is often a small RAM-backed tmpfs; point it at
# real disk since uploads spool there before the app sees them.
Environment=TMPDIR=/opt/shred/tmp

# gthread lets a worker serve other requests during a slow upload/download
# instead of blocking the whole process. Tune --workers/--threads to your
# own hardware — these are starting points, not a requirement.
ExecStart=/opt/shred/venv/bin/gunicorn --workers 2 --worker-class gthread --threads 4 --bind 127.0.0.1:8000 --timeout 3600 server:app

Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/shred/data /opt/shred/tmp
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=false
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
SystemCallFilter=@system-service
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target